# package/panicos-input-sense/files/system.d/hdmi-hotplug.path
[Unit]
Description=Monitor HDMI status
StartLimitIntervalSec=30
StartLimitBurst=20
StartLimitAction=none

[Path]
PathModified=/run/hdmi-status.last

[Install]
WantedBy=multi-user.target
