# package/panicos-input-sense/files/system.d/hdmi-hotplug.service
[Unit]
Description=Handle HDMI hotplug
# A single HDMI plug/unplug generates many rapid DRM change events.
# Disable rate limiting entirely so the path unit never enters failed state.
StartLimitIntervalSec=0

[Service]
ExecStartPre=/bin/sleep 1
ExecStart=/usr/bin/handle-hdmi-hotplug
Type=oneshot
