# package/panicos-pht/panicos-pht.service
[Unit]
Description=ProHandheldTracker
# Order after the basics + audio + network so PHT has everything it needs
# at startup, but DON'T put After=multi-user.target — that creates a
# subtle ordering issue with WantedBy=multi-user.target where systemd
# silently drops the activation. After basic.target is the standard
# pattern for multi-user services.
After=basic.target sound.target network.target

[Service]
Type=simple
ExecStart=/opt/pht/panictracker.sh
# No TTYPath / StandardInput=tty here on purpose. PHT renders via KMSDRM
# (grabs DRM master directly on /dev/dri/card0); it doesn't need TTY
# semantics. Keeping a TTY claim was getting systemd into TTY-ownership
# loops with the splash service.
StandardInput=null
StandardOutput=journal
StandardError=journal
# Restart on crash so the device returns to PHT, not a black screen.
Restart=on-failure
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
