# package/panicos-sshkeys/panicos-sshkeys.service
[Unit]
Description=PanicOS first-boot: generate per-device SSH host keys
DefaultDependencies=no
After=local-fs.target
Before=basic.target sysinit.target dropbear.service
ConditionPathExists=!/storage/.panicos-sshkeys-done
RequiresMountsFor=/storage

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/panicos-sshkeys
StandardOutput=journal+console

[Install]
WantedBy=sysinit.target
