# flavors/launcher/rootfs-overlay/usr/lib/systemd/system/wireplumber.service.d/panicos-kiosk.conf
# PanicOS launcher kiosk drop-in for wireplumber.service.
# Buildroot's wireplumber.service is meant for per-user systemd
# instances; we run as root with no user session. The drop-in:
#   - HOME / XDG_STATE_HOME / runtime/state dirs so wireplumber's
#     .local/state writes go somewhere writable instead of crashing
#     with "failed to create directory /.local/state/wireplumber:
#     Permission denied"
#   - DBUS_SESSION_BUS_ADDRESS pointing at the system bus so the
#     bluez5 monitor module (and other D-Bus consumers) can talk to
#     bluetoothd. Without this, BT audio device hotplug never reaches
#     wireplumber and BT speakers/headphones don't appear as sinks.
[Service]
Environment=HOME=/var/lib/wireplumber
Environment=XDG_STATE_HOME=/var/lib/wireplumber
Environment=XDG_RUNTIME_DIR=/run
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/dbus/system_bus_socket
RuntimeDirectory=wireplumber
StateDirectory=wireplumber
