# flavors/launcher/rootfs-overlay/usr/lib/systemd/system/wireplumber.service.d/panicos-kiosk.conf
# PanicOS launcher kiosk drop-in for wireplumber.service.
#
# Mirrors ROCKNIX's audio/wireplumber/patches/001-fix-service.patch
# changes to the upstream wireplumber.service. Buildroot ships
# wireplumber's stock unit which runs as User=pipewire — that user
# isn't in the audio group, so wireplumber starts but its ALSA monitor
# enumerates zero cards (`wpctl status` shows empty Sinks/Sources and
# Default Configured Devices = (none)). Switching to User=root
# sidesteps the missing-group-membership problem.
#
# Also pins the system bus / runtime dir / writable HOME for
# wireplumber's .local/state writes. Together with the conf.d disables
# under /usr/share/wireplumber/wireplumber.conf.d/ (89-disable-*),
# this is the same combination ROCKNIX uses for system-mode wireplumber.
[Service]
User=root
Environment=HOME=/storage
Environment=PIPEWIRE_RUNTIME_DIR=/run/pipewire
Environment=XDG_RUNTIME_DIR=/var/run/0-runtime-dir
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/dbus/system_bus_socket
RuntimeDirectory=wireplumber
StateDirectory=wireplumber
