# package/panicos-emulationstation/files/panicos-es.service
[Unit]
Description=EmulationStation launcher (Wayland client of panicos-sway)
# Sway provides the panel; ES connects to its Wayland socket.
# Mirrors ROCKNIX's essway.service design.
After=panicos-sway.service basic.target sound.target network.target
Requires=panicos-sway.service

[Service]
Type=simple
# Wayland environment so SDL2 picks the wayland video backend rather
# than KMSDRM. Sway puts its socket at $XDG_RUNTIME_DIR/wayland-1.
Environment=HOME=/storage
Environment=XDG_RUNTIME_DIR=/var/run/0-runtime-dir
Environment=WAYLAND_DISPLAY=wayland-1
Environment=SDL_VIDEODRIVER=wayland
ExecStart=/usr/bin/emulationstation
StandardInput=null
StandardOutput=journal
StandardError=journal
# ES sometimes exits cleanly when the user picks Quit; restart so the
# device returns to the launcher rather than a black screen.
Restart=always
RestartSec=2
User=root
WorkingDirectory=/storage

[Install]
WantedBy=multi-user.target
