# package/panicos-emulationstation/files/panicos-es.service
[Unit]
Description=EmulationStation launcher (Wayland client of panicos-sway)
# Sway provides the panel; ES connects to its Wayland socket.
# Mirrors ROCKNIX's essway.service design.
After=panicos-sway.service basic.target sound.target network.target
Requires=panicos-sway.service

[Service]
Type=simple
# Wayland environment so SDL2 picks the wayland video backend rather
# than KMSDRM. Sway puts its socket at $XDG_RUNTIME_DIR/wayland-1.
Environment=HOME=/storage
Environment=XDG_RUNTIME_DIR=/var/run/0-runtime-dir
Environment=WAYLAND_DISPLAY=wayland-1
Environment=SDL_VIDEODRIVER=wayland
# SWAYSOCK is needed by every port that calls swaymsg (PortMaster's
# pm_platform_helper, port wrappers that fullscreen themselves).
# sway-100.01-static-ipc-socket.patch makes the path predictable (no PID).
Environment=SWAYSOCK=/var/run/0-runtime-dir/sway-ipc.0.sock
# OpenAL Soft default backend list picks alsa first, which routes
# through the pipewire-alsa shim and produces "Wait timeout... buffer
# size too low?" spam under load (see Doom Engines / Wolfenstein log).
# Native pipewire backend avoids the shim. Inherited by every port.
Environment=ALSOFT_DRIVERS=pipewire,alsa,sdl2
# Force SDL audio through PulseAudio (handled by PipeWire's compat
# server). Without this SDL tries ALSA first, bypassing PipeWire's
# sink routing — Rockbox and other SDL apps then miss headphone
# hotplug events and stay routed to the speaker even when USB-C
# headphones are plugged in. Inherited by all PortMaster ports.
Environment=SDL_AUDIODRIVER=pulseaudio
# Wait up to 10 s for sway to create its Wayland socket before launching ES.
# sway 1.11 has no sd_notify, so Type=simple on the sway unit means systemd
# marks it "active" the instant the process forks — before the socket exists.
ExecStartPre=/bin/sh -c 'i=0; while [ $i -lt 100 ]; do [ -S "${XDG_RUNTIME_DIR}/${WAYLAND_DISPLAY}" ] && exit 0; sleep 0.1; i=$((i+1)); done; echo "panicos-es: timed out waiting for wayland socket" >&2; exit 1'
ExecStart=/usr/bin/emulationstation
StandardInput=null
StandardOutput=journal
StandardError=journal
# ES sometimes exits cleanly when the user picks Quit; restart so the
# device returns to the launcher rather than a black screen.
Restart=always
RestartSec=2
User=root
WorkingDirectory=/storage

[Install]
WantedBy=multi-user.target
