# package/panicos-firstboot/panicos-firstboot.service
[Unit]
Description=PanicOS first-boot: grow overlay partition
DefaultDependencies=no
After=local-fs.target
Before=basic.target sysinit.target
ConditionPathExists=!/storage/.panicos-firstboot-done
RequiresMountsFor=/storage

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/panicos-firstboot
StandardOutput=journal+console

[Install]
WantedBy=sysinit.target
