# package/panicos-firstboot/panicos-firstboot.service
[Unit]
Description=PanicOS first-boot: grow overlay partition
DefaultDependencies=no
After=local-fs.target
Before=basic.target sysinit.target
ConditionPathExists=!/storage/.panicos-firstboot-done
RequiresMountsFor=/storage

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/sh -c 'echo -ne "\033[1000H\033[2K==> First-boot setup, please wait..." >/dev/console'
ExecStart=/usr/sbin/panicos-firstboot
StandardOutput=journal

[Install]
WantedBy=sysinit.target
