# package/panicos-wifi-config/panicos-wifi-config.service
[Unit]
Description=PanicOS: render wpa_supplicant config from boot vfat
DefaultDependencies=no
After=local-fs.target
Before=basic.target wpa_supplicant@wlan0.service
RequiresMountsFor=/boot

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/panicos-wifi-config
StandardOutput=journal

[Install]
WantedBy=sysinit.target
