# package/panicos-wifi-config/panicos-wifi-config.service
[Unit]
Description=PanicOS: render wpa_supplicant config from boot vfat
DefaultDependencies=no
After=local-fs.target
Before=basic.target wpa_supplicant@wlan0.service
RequiresMountsFor=/boot

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/sh -c 'echo -ne "\033[1000H\033[2K==> Applying Wi-Fi config..." >/dev/console'
ExecStart=/usr/sbin/panicos-wifi-config
StandardOutput=journal

[Install]
WantedBy=sysinit.target
