# flavors/launcher/rootfs-overlay/usr/lib/systemd/system/pipewire.service.d/panicos-kiosk.conf
# PanicOS launcher kiosk drop-in for pipewire.service.
#
# Mirrors ROCKNIX's audio/pipewire/patches/001-pipewire-pulse.patch
# changes to the upstream pipewire.service. Buildroot ships pipewire's
# stock unit which runs as User=pipewire — that user exists but isn't
# in the audio group, so the daemon can't access /dev/snd/*. ROCKNIX
# bumps it to root and pins the system bus / runtime dir env so the
# daemon and its bluez5 module can reach bluetoothd.
[Service]
User=root
Environment=HOME=/var/lib/pipewire
Environment=XDG_RUNTIME_DIR=/var/run/0-runtime-dir
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/dbus/system_bus_socket
StateDirectory=pipewire
