# flavors/launcher/rootfs-overlay/usr/lib/systemd/system/pipewire.service.d/panicos-kiosk.conf
# PanicOS launcher kiosk drop-in for pipewire.service.
# Same rationale as the wireplumber drop-in: kiosk runs as root with
# no user session. Aliasing DBUS_SESSION_BUS_ADDRESS to the system bus
# socket lets pipewire's bluez5 module reach bluetoothd, which is the
# path BT audio devices arrive through (bluetoothd → system bus →
# pipewire bluez5 module → wireplumber sink → ALSA via pipewire-alsa).
[Service]
Environment=HOME=/var/lib/pipewire
Environment=XDG_RUNTIME_DIR=/run
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/dbus/system_bus_socket
StateDirectory=pipewire
