# package/panicos-splash/panicos-splash.service
[Unit]
Description=PanicOS boot splash
DefaultDependencies=no
After=local-fs.target
Before=basic.target getty@tty1.service panicos-pht.service
# Stop the splash the moment the real UI starts. systemd handles the
# transition cleanly: SIGTERM us → start the conflicting unit.
Conflicts=getty@tty1.service panicos-pht.service

[Service]
Type=simple
ExecStart=/usr/sbin/panicos-splash
# Don't claim tty1 via TTYPath — fbv writes directly to /dev/fb0 (which
# any root process can do) and doesn't need TTY semantics. Holding tty1
# here means systemd refuses to schedule panicos-pht.service (which DOES
# need tty1) until we release it, which we never do voluntarily — so pht
# stayed in "inactive (dead)" state forever waiting on us.
StandardInput=null
StandardOutput=journal
StandardError=journal
KillMode=process
TimeoutStopSec=2
# fbv never exits voluntarily (no --delay); systemd needs to know that's
# OK and not flag it as a failure when we're killed during transition.
SuccessExitStatus=SIGTERM SIGKILL
Restart=no

[Install]
WantedBy=sysinit.target
