# package/panicos-splash/panicos-splash.service
[Unit]
Description=PanicOS boot splash
DefaultDependencies=no
After=local-fs.target
Before=basic.target

[Service]
# Type=oneshot — fbv (or whatever the launcher uses) draws the image
# and exits; the image stays painted on the framebuffer until something
# else (PHT via KMSDRM, getty via fbcon) overwrites it. No TTY claim, no
# Conflicts= — those were causing systemd to refuse to start panicos-pht.
Type=oneshot
ExecStart=/usr/sbin/panicos-splash
StandardInput=null
StandardOutput=journal
StandardError=journal
RemainAfterExit=no

[Install]
WantedBy=sysinit.target
